# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:riemann]
command=/riemann/bin/riemann
autorestart=true

[program:dashboard]
command=/usr/local/bin/riemann-dash
autorestart=true

[program:varnishd]
command=/usr/sbin/varnishd -s malloc,1G -T 127.0.0.1:6082 -a 0.0.0.0:6081 -F -f /docker-riemann/riemann.vcl
autorestart=true
