# sv/riemann/run
#!/bin/sh

#sv start socklog-unix || exit 1

exec 2>&1
exec /riemann/bin/riemann
